# examples/systemd/pistis-stage.path
[Unit]
Description=Git repository monitor

[Path]
PathChanged=/srv/git-stage/.git/refs/heads/main
TriggerLimitBurst=50

[Install]
WantedBy=paths.target
